# deploy/launchd/interval-report.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__APP_DIR__/run_report.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__APP_DIR__</string>
  <key>StartInterval</key>
  <integer>__SECONDS__</integer>
  <key>StandardOutPath</key>
  <string>__APP_DIR__/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__APP_DIR__/logs/launchd.err.log</string>
</dict>
</plist>
